# etc/systemd/system/codex-setup.service
[Unit]
Description=Heal and run setup.sh before networking
DefaultDependencies=no
Before=network-pre.target
Wants=network-pre.target

[Service]
Type=oneshot
ExecStart=/usr/local/bin/codex -q -a full-auto "doctor setup.sh" && /usr/local/bin/setup.sh
RemainAfterExit=yes

[Install]
WantedBy=multi-user.target
